# forked-daapd/rootfs/etc/services.d/forked-daapd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Ulrar Hass.io Add-ons: forked-daapd server
# Runs the forked-daapd Server
# ==============================================================================
# shellcheck disable=SC1091

bashio::log.info 'Starting the forked-daapd Server...'
# Run forked-daapd
exec forked-daapd -f -c /share/forked-daapd/forked-daapd.conf -P /var/run/forked-daapd.pid
