# forked-daapd/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Ulrar Hass.io Add-ons: forked-daapd server
# Runs the forked-daapd Server
# ==============================================================================
# shellcheck disable=SC1091
declare -a options

# Wait for Forked-Daapd to become available
bashio::net.wait_for 3689

bashio::log.info "Starting Nginx..."

exec nginx "${options[@]}"
