# radicale/rootfs/etc/services.d/radicale/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Radicale addon
# ==============================================================================
declare -a options

bashio::log.info "Starting Radicale ..."

radicale --config /etc/radicale/config
